# crontab
23 23 * * * /Users/jinkwon/Desktop/Python/Hacker News Headling Emailer/dist/news scraper
